# node.module.club.conf
# node.module.club
# Generated by strong-service-upstart on 2014-08-09T19:11:55.149Z

description "A node app"
author "root"

# Start on boot/reboot
start on runlevel [2345]

# Stop on shutdown/reboot
stop on runlevel [016]

# To start/stop node.module.club service manually, run:
#
# /sbin/initctl start node.module.club
# /sbin/initctl stop node.module.club

respawn

# Log all service output to /var/log/node.module.club.log
# Log rotation is handled by Upstart
console log

# User & Group to run service as:
chdir /root/site
setuid root 
setgid nogroup

# Arbitrarily high limit on number of open file descriptors:
limit nofile 50000 50000

# Enable core dumps:
limit core unlimited unlimited

# Environment variables:
# no environment variables set

# Run our actual application
exec /usr/local/bin/swurve -h module.club
